# docker_conf/runit/cron/run
#!/bin/sh
exec crond -f -c /app/docker_conf/crontab
